# docker/app/supervisord.conf
[supervisord]
# The application container is intentionally multi-process; Supervisor keeps the
# daemons in the foreground and restarts them independently if needed.
nodaemon=true

[program:avahi]
command=/usr/local/bin/run-avahi.sh
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mpd]
# Start order matters: MPD comes up before myMPD and mpdscribble connect to it.
command=/usr/local/bin/run-mpd.sh
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mympd]
# myMPD depends on the local MPD instance and serves the main web interface.
command=/usr/local/bin/run-mympd.sh
autorestart=true
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:minidlna]
# MiniDLNA is independent from MPD but shares the same mounted music library.
command=/usr/local/bin/run-minidlna.sh
autorestart=true
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:snapserver]
# Snapserver consumes the FIFO audio output emitted by MPD when enabled.
command=/usr/local/bin/run-snapserver.sh
autorestart=true
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:mpdscribble]
# Scrobbling starts late because it requires MPD plus user credentials.
command=/usr/local/bin/run-mpdscribble.sh
autorestart=true
priority=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ddns_cron]
# The FreeDNS updater is optional and simply idles when UPDATE_URL is blank.
command=/usr/local/bin/run-ddns-cron.sh
autorestart=true
priority=70
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
